# Linux_service_sensor/systemd/assignment-sensor.service
[Unit]
Description=Mock sensor logger (assignment)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/assignment-sensor --interval 2 --logfile /var/tmp/sensor.log
Restart=always
RestartSec=2
User=nobody
WorkingDirectory=/tmp

[Install]
WantedBy=multi-user.target
